# udev/80-closeremove.rules
ACTION=="remove", ATTRS{idProduct}=="0010", RUN+="/bin/loginctl lock-sessions"
